# data/fixtures/attack_data_t1053_sysmon.json
[
  {
    "EventTime": "2020-08-01 06:49:15",
    "Hostname": "win-dc-8400769.attackrange.local",
    "Keywords": -9223372036854775808,
    "EventType": "INFO",
    "SeverityValue": 2,
    "Severity": "INFO",
    "EventID": 10,
    "SourceName": "Microsoft-Windows-Sysmon",
    "ProviderGuid": "{5770385F-C22A-43E0-BF4C-06F5698FFBD9}",
    "Version": 3,
    "Task": 10,
    "OpcodeValue": 0,
    "RecordNumber": 838,
    "ProcessID": 2928,
    "ThreadID": 3404,
    "Channel": "Microsoft-Windows-Sysmon/Operational",
    "Domain": "NT AUTHORITY",
    "AccountName": "SYSTEM",
    "UserID": "S-1-5-18",
    "AccountType": "User",
    "Message": "Process accessed:\r\nRuleName: \r\nUtcTime: 2020-08-01 06:48:59.225\r\nSourceProcessGUID: {E2A3D6B1-105A-5F25-0000-0010CE650000}\r\nSourceProcessId: 596\r\nSourceThreadId: 728\r\nSourceImage: C:\\Windows\\system32\\svchost.exe\r\nTargetProcessGUID: {E2A3D6B1-1058-5F25-0000-001012540000}\r\nTargetProcessId: 864\r\nTargetImage: C:\\Windows\\system32\\lsass.exe\r\nGrantedAccess: 0x1400\r\nCallTrace: C:\\Windows\\SYSTEM32\\ntdll.dll+a5ec4|C:\\Windows\\System32\\KERNELBASE.dll+5eab4|c:\\windows\\system32\\lsm.dll+fa3b|C:\\Windows\\System32\\RPCRT4.dll+78253|C:\\Windows\\System32\\RPCRT4.dll+dbc0d|C:\\Windows\\System32\\RPCRT4.dll+b3cc|C:\\Windows\\System32\\RPCRT4.dll+59cd4|C:\\Windows\\System32\\RPCRT4.dll+58bed|C:\\Windows\\System32\\RPCRT4.dll+5949b|C:\\Windows\\System32\\RPCRT4.dll+3933c|C:\\Windows\\System32\\RPCRT4.dll+397bc|C:\\Windows\\System32\\RPCRT4.dll+53dac|C:\\Windows\\System32\\RPCRT4.dll+5560b|C:\\Windows\\System32\\RPCRT4.dll+480ea|C:\\Windows\\SYSTEM32\\ntdll.dll+286be|C:\\Windows\\SYSTEM32\\ntdll.dll+2a029|C:\\Windows\\System32\\KERNEL32.DLL+84d4|C:\\Windows\\SYSTEM32\\ntdll.dll+6e871",
    "Category": "Process accessed (rule: ProcessAccess)",
    "Opcode": "Info",
    "UtcTime": "2020-08-01 06:48:59.225",
    "SourceProcessGUID": "{E2A3D6B1-105A-5F25-0000-0010CE650000}",
    "SourceProcessId": "596",
    "SourceThreadId": "728",
    "SourceImage": "C:\\Windows\\system32\\svchost.exe",
    "TargetProcessGUID": "{E2A3D6B1-1058-5F25-0000-001012540000}",
    "TargetProcessId": "864",
    "TargetImage": "C:\\Windows\\system32\\lsass.exe",
    "GrantedAccess": "0x1400",
    "CallTrace": "C:\\Windows\\SYSTEM32\\ntdll.dll+a5ec4|C:\\Windows\\System32\\KERNELBASE.dll+5eab4|c:\\windows\\system32\\lsm.dll+fa3b|C:\\Windows\\System32\\RPCRT4.dll+78253|C:\\Windows\\System32\\RPCRT4.dll+dbc0d|C:\\Windows\\System32\\RPCRT4.dll+b3cc|C:\\Windows\\System32\\RPCRT4.dll+59cd4|C:\\Windows\\System32\\RPCRT4.dll+58bed|C:\\Windows\\System32\\RPCRT4.dll+5949b|C:\\Windows\\System32\\RPCRT4.dll+3933c|C:\\Windows\\System32\\RPCRT4.dll+397bc|C:\\Windows\\System32\\RPCRT4.dll+53dac|C:\\Windows\\System32\\RPCRT4.dll+5560b|C:\\Windows\\System32\\RPCRT4.dll+480ea|C:\\Windows\\SYSTEM32\\ntdll.dll+286be|C:\\Windows\\SYSTEM32\\ntdll.dll+2a029|C:\\Windows\\System32\\KERNEL32.DLL+84d4|C:\\Windows\\SYSTEM32\\ntdll.dll+6e871",
    "EventReceivedTime": "2020-08-01 06:49:16",
    "SourceModuleName": "eventlog",
    "SourceModuleType": "im_msvistalog"
  },
  {
    "EventTime": "2020-08-01 06:49:15",
    "Hostname": "win-dc-8400769.attackrange.local",
    "Keywords": -9223372036854775808,
    "EventType": "INFO",
    "SeverityValue": 2,
    "Severity": "INFO",
    "EventID": 13,
    "SourceName": "Microsoft-Windows-Sysmon",
    "ProviderGuid": "{5770385F-C22A-43E0-BF4C-06F5698FFBD9}",
    "Version": 2,
    "Task": 13,
    "OpcodeValue": 0,
    "RecordNumber": 903,
    "ProcessID": 2928,
    "ThreadID": 3404,
    "Channel": "Microsoft-Windows-Sysmon/Operational",
    "Domain": "NT AUTHORITY",
    "AccountName": "SYSTEM",
    "UserID": "S-1-5-18",
    "AccountType": "User",
    "Message": "Registry value set:\r\nRuleName: \r\nEventType: SetValue\r\nUtcTime: 2020-08-01 06:48:59.350\r\nProcessGuid: {E2A3D6B1-105B-5F25-0000-001015C40000}\r\nProcessId: 1196\r\nImage: C:\\Windows\\system32\\svchost.exe\r\nTargetObject: HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\NetworkList\\Profiles\\{1C4380BB-3862-41F1-B425-7E74397CF5D1}\\DateLastConnected\r\nDetails: Binary Data",
    "Category": "Registry value set (rule: RegistryEvent)",
    "Opcode": "Info",
    "UtcTime": "2020-08-01 06:48:59.350",
    "ProcessGuid": "{E2A3D6B1-105B-5F25-0000-001015C40000}",
    "Image": "C:\\Windows\\system32\\svchost.exe",
    "TargetObject": "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\NetworkList\\Profiles\\{1C4380BB-3862-41F1-B425-7E74397CF5D1}\\DateLastConnected",
    "Details": "Binary Data",
    "EventReceivedTime": "2020-08-01 06:49:16",
    "SourceModuleName": "eventlog",
    "SourceModuleType": "im_msvistalog"
  },
  {
    "EventTime": "2020-08-01 06:49:15",
    "Hostname": "win-dc-8400769.attackrange.local",
    "Keywords": -9223372036854775808,
    "EventType": "INFO",
    "SeverityValue": 2,
    "Severity": "INFO",
    "EventID": 11,
    "SourceName": "Microsoft-Windows-Sysmon",
    "ProviderGuid": "{5770385F-C22A-43E0-BF4C-06F5698FFBD9}",
    "Version": 2,
    "Task": 11,
    "OpcodeValue": 0,
    "RecordNumber": 905,
    "ProcessID": 2928,
    "ThreadID": 3404,
    "Channel": "Microsoft-Windows-Sysmon/Operational",
    "Domain": "NT AUTHORITY",
    "AccountName": "SYSTEM",
    "UserID": "S-1-5-18",
    "AccountType": "User",
    "Message": "File created:\r\nRuleName: T1053\r\nUtcTime: 2020-08-01 06:48:59.381\r\nProcessGuid: {E2A3D6B1-105B-5F25-0000-001056BA0000}\r\nProcessId: 1140\r\nImage: C:\\Windows\\system32\\svchost.exe\r\nTargetFilename: C:\\Windows\\Tasks\\SA.DAT\r\nCreationUtcTime: 2016-09-12 11:34:03.403",
    "Category": "File created (rule: FileCreate)",
    "Opcode": "Info",
    "RuleName": "T1053",
    "UtcTime": "2020-08-01 06:48:59.381",
    "ProcessGuid": "{E2A3D6B1-105B-5F25-0000-001056BA0000}",
    "Image": "C:\\Windows\\system32\\svchost.exe",
    "TargetFilename": "C:\\Windows\\Tasks\\SA.DAT",
    "CreationUtcTime": "2016-09-12 11:34:03.403",
    "EventReceivedTime": "2020-08-01 06:49:16",
    "SourceModuleName": "eventlog",
    "SourceModuleType": "im_msvistalog"
  },
  {
    "EventTime": "2020-08-01 06:49:18",
    "Hostname": "win-dc-8400769.attackrange.local",
    "Keywords": -9223372036854775808,
    "EventType": "INFO",
    "SeverityValue": 2,
    "Severity": "INFO",
    "EventID": 3,
    "SourceName": "Microsoft-Windows-Sysmon",
    "ProviderGuid": "{5770385F-C22A-43E0-BF4C-06F5698FFBD9}",
    "Version": 5,
    "Task": 3,
    "OpcodeValue": 0,
    "RecordNumber": 1996,
    "ProcessID": 2928,
    "ThreadID": 3388,
    "Channel": "Microsoft-Windows-Sysmon/Operational",
    "Domain": "NT AUTHORITY",
    "AccountName": "SYSTEM",
    "UserID": "S-1-5-18",
    "AccountType": "User",
    "Message": "Network connection detected:\r\nRuleName: Usermode\r\nUtcTime: 2020-08-01 06:49:16.509\r\nProcessGuid: {E2A3D6B1-1061-5F25-0000-001069980200}\r\nProcessId: 2308\r\nImage: C:\\Users\\Public\\splunkd.exe\r\nUser: NT AUTHORITY\\SYSTEM\r\nProtocol: tcp\r\nInitiated: true\r\nSourceIsIpv6: false\r\nSourceIp: 10.0.1.14\r\nSourceHostname: win-dc-8400769.attackrange.local\r\nSourcePort: 49689\r\nSourcePortName: \r\nDestinationIsIpv6: false\r\nDestinationIp: 10.0.1.12\r\nDestinationHostname: \r\nDestinationPort: 7010\r\nDestinationPortName: ",
    "Category": "Network connection detected (rule: NetworkConnect)",
    "Opcode": "Info",
    "RuleName": "Usermode",
    "UtcTime": "2020-08-01 06:49:16.509",
    "ProcessGuid": "{E2A3D6B1-1061-5F25-0000-001069980200}",
    "Image": "C:\\Users\\Public\\splunkd.exe",
    "User": "NT AUTHORITY\\SYSTEM",
    "Protocol": "tcp",
    "Initiated": "true",
    "SourceIsIpv6": "false",
    "SourceIp": "10.0.1.14",
    "SourceHostname": "win-dc-8400769.attackrange.local",
    "SourcePort": "49689",
    "DestinationIsIpv6": "false",
    "DestinationIp": "10.0.1.12",
    "DestinationPort": "7010",
    "EventReceivedTime": "2020-08-01 06:49:19",
    "SourceModuleName": "eventlog",
    "SourceModuleType": "im_msvistalog"
  },
  {
    "EventTime": "2020-08-01 06:55:19",
    "Hostname": "win-dc-8400769.attackrange.local",
    "Keywords": -9223372036854775808,
    "EventType": "INFO",
    "SeverityValue": 2,
    "Severity": "INFO",
    "EventID": 1,
    "SourceName": "Microsoft-Windows-Sysmon",
    "ProviderGuid": "{5770385F-C22A-43E0-BF4C-06F5698FFBD9}",
    "Version": 5,
    "Task": 1,
    "OpcodeValue": 0,
    "RecordNumber": 5674,
    "ProcessID": 2928,
    "ThreadID": 3404,
    "Channel": "Microsoft-Windows-Sysmon/Operational",
    "Domain": "NT AUTHORITY",
    "AccountName": "SYSTEM",
    "UserID": "S-1-5-18",
    "AccountType": "User",
    "Message": "Process Create:\r\nRuleName: \r\nUtcTime: 2020-08-01 06:55:19.116\r\nProcessGuid: {E2A3D6B1-11D7-5F25-0000-001068550800}\r\nProcessId: 4924\r\nImage: C:\\Windows\\System32\\cmd.exe\r\nFileVersion: 10.0.14393.0 (rs1_release.160715-1616)\r\nDescription: Windows Command Processor\r\nProduct: Microsoft\u00ae Windows\u00ae Operating System\r\nCompany: Microsoft Corporation\r\nOriginalFileName: Cmd.Exe\r\nCommandLine: \"C:\\Windows\\system32\\cmd.exe\" /c \"schtasks /create /tn \"T1053_005_OnLogon\" /sc onlogon /tr \"cmd.exe /c calc.exe\" & schtasks /create /tn \"T1053_005_OnStartup\" /sc onstart /ru system /tr \"cmd.exe /c calc.exe\"\" \r\nCurrentDirectory: C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\\r\nUser: ATTACKRANGE\\Administrator\r\nLogonGuid: {E2A3D6B1-11D4-5F25-0000-00205F0B0800}\r\nLogonId: 0x80B5F\r\nTerminalSessionId: 0\r\nIntegrityLevel: High\r\nHashes: MD5=F4F684066175B77E0C3A000549D2922C,SHA256=935C1861DF1F4018D698E8B65ABFA02D7E9037D8F68CA3C2065B6CA165D44AD2,IMPHASH=3062ED732D4B25D1C64F084DAC97D37A\r\nParentProcessGuid: {E2A3D6B1-11D6-5F25-0000-00106A320800}\r\nParentProcessId: 4848\r\nParentImage: C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe\r\nParentCommandLine: \"powershell.exe\" -noninteractive -encodedcommand WwBDAG8AbgBzAG8AbABlAF0AOgA6AEkAbgBwAHUAdABFAG4AYwBvAGQAaQBuAGcAIAA9ACAATgBlAHcALQBPAGIAagBlAGMAdAAgAFQAZQB4AHQALgBVAFQARgA4AEUAbgBjAG8AZABpAG4AZwAgACQAZgBhAGwAcwBlADsAIABJAG0AcABvAHIAdAAtAE0AbwBkAHUAbABlACAAIgBDADoAXABBAHQAbwBtAGkAYwBSAGUAZABUAGUAYQBtAFwAaQBuAHYAbwBrAGUALQBhAHQAbwBtAGkAYwByAGUAZAB0AGUAYQBtAFwASQBuAHYAbwBrAGUALQBBAHQAbwBtAGkAYwBSAGUAZABUAGUAYQBtAC4AcABzAGQAMQAiACAALQBGAG8AcgBjAGUACgAkAG0AeQBBAHIAZwBzACAAPQAgAEAAewAgACIAdQBzAGUAcgBfAG4AYQBtAGUAIgAgAD0AIAAiAEEAVABUAEEAQwBLAFIAQQBOAEcARQBcAEEAZABtAGkAbgBpAHMAdAByAGEAdABvAHIAIgA7ACAAIgBwAGEAcwBzAHcAbwByAGQAIgAgAD0AIAAiAEkALQBsADEAawBlAC0AQQB0AHQAYQBjAGsALQBSAGEAbgBnAGUAIQAiACAAfQAKAEkAbgB2AG8AawBlAC0AQQB0AG8AbQBpAGMAVABlAHMAdAAgACIAVAAxADAANQAzAC4AMAAwADUAIgAgAC0ASQBuAHAAdQB0AEEAcgBnAHMAIAAkAG0AeQBBAHIAZwBzACAALQBDAG8AbgBmAGkAcgBtADoAJABmAGEAbABzAGUAIAAtAFQAaQBtAGUAbwB1AHQAUwBlAGMAbwBuAGQAcwAgADMAMAAwACAALQBFAHgAZQBjAHUAdABpAG8AbgBMAG8AZwBQAGEAdABoACAAQwA6AFwAQQB0AG8AbQBpAGMAUgBlAGQAVABlAGEAbQBcAGEAdABjAF8AZQB4AGUAYwB1AHQAaQBvAG4ALgBjAHMAdgA=",
    "Category": "Process Create (rule: ProcessCreate)",
    "Opcode": "Info",
    "UtcTime": "2020-08-01 06:55:19.116",
    "ProcessGuid": "{E2A3D6B1-11D7-5F25-0000-001068550800}",
    "Image": "C:\\Windows\\System32\\cmd.exe",
    "FileVersion": "10.0.14393.0 (rs1_release.160715-1616)",
    "Description": "Windows Command Processor",
    "Product": "Microsoft\u00ae Windows\u00ae Operating System",
    "Company": "Microsoft Corporation",
    "OriginalFileName": "Cmd.Exe",
    "CommandLine": "\"C:\\Windows\\system32\\cmd.exe\" /c \"schtasks /create /tn \"T1053_005_OnLogon\" /sc onlogon /tr \"cmd.exe /c calc.exe\" &amp; schtasks /create /tn \"T1053_005_OnStartup\" /sc onstart /ru system /tr \"cmd.exe /c calc.exe\"\" ",
    "CurrentDirectory": "C:\\Users\\ADMINI~1\\AppData\\Local\\Temp\\",
    "User": "ATTACKRANGE\\Administrator",
    "LogonGuid": "{E2A3D6B1-11D4-5F25-0000-00205F0B0800}",
    "LogonId": "0x80b5f",
    "TerminalSessionId": "0",
    "IntegrityLevel": "High",
    "Hashes": "MD5=F4F684066175B77E0C3A000549D2922C,SHA256=935C1861DF1F4018D698E8B65ABFA02D7E9037D8F68CA3C2065B6CA165D44AD2,IMPHASH=3062ED732D4B25D1C64F084DAC97D37A",
    "ParentProcessGuid": "{E2A3D6B1-11D6-5F25-0000-00106A320800}",
    "ParentProcessId": "4848",
    "ParentImage": "C:\\Windows\\System32\\WindowsPowerShell\\v1.0\\powershell.exe",
    "ParentCommandLine": "\"powershell.exe\" -noninteractive -encodedcommand WwBDAG8AbgBzAG8AbABlAF0AOgA6AEkAbgBwAHUAdABFAG4AYwBvAGQAaQBuAGcAIAA9ACAATgBlAHcALQBPAGIAagBlAGMAdAAgAFQAZQB4AHQALgBVAFQARgA4AEUAbgBjAG8AZABpAG4AZwAgACQAZgBhAGwAcwBlADsAIABJAG0AcABvAHIAdAAtAE0AbwBkAHUAbABlACAAIgBDADoAXABBAHQAbwBtAGkAYwBSAGUAZABUAGUAYQBtAFwAaQBuAHYAbwBrAGUALQBhAHQAbwBtAGkAYwByAGUAZAB0AGUAYQBtAFwASQBuAHYAbwBrAGUALQBBAHQAbwBtAGkAYwBSAGUAZABUAGUAYQBtAC4AcABzAGQAMQAiACAALQBGAG8AcgBjAGUACgAkAG0AeQBBAHIAZwBzACAAPQAgAEAAewAgACIAdQBzAGUAcgBfAG4AYQBtAGUAIgAgAD0AIAAiAEEAVABUAEEAQwBLAFIAQQBOAEcARQBcAEEAZABtAGkAbgBpAHMAdAByAGEAdABvAHIAIgA7ACAAIgBwAGEAcwBzAHcAbwByAGQAIgAgAD0AIAAiAEkALQBsADEAawBlAC0AQQB0AHQAYQBjAGsALQBSAGEAbgBnAGUAIQAiACAAfQAKAEkAbgB2AG8AawBlAC0AQQB0AG8AbQBpAGMAVABlAHMAdAAgACIAVAAxADAANQAzAC4AMAAwADUAIgAgAC0ASQBuAHAAdQB0AEEAcgBnAHMAIAAkAG0AeQBBAHIAZwBzACAALQBDAG8AbgBmAGkAcgBtADoAJABmAGEAbABzAGUAIAAtAFQAaQBtAGUAbwB1AHQAUwBlAGMAbwBuAGQAcwAgADMAMAAwACAALQBFAHgAZQBjAHUAdABpAG8AbgBMAG8AZwBQAGEAdABoACAAQwA6AFwAQQB0AG8AbQBpAGMAUgBlAGQAVABlAGEAbQBcAGEAdABjAF8AZQB4AGUAYwB1AHQAaQBvAG4ALgBjAHMAdgA=",
    "EventReceivedTime": "2020-08-01 06:55:20",
    "SourceModuleName": "eventlog",
    "SourceModuleType": "im_msvistalog"
  }
]